# ranch-camera.service
[Unit]
Description=Ranch Camera Capture Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pi
WorkingDirectory=/home/pi
ExecStart=/usr/bin/python3 /home/pi/capture_upload.py

# Supabase Configuration
Environment="SUPABASE_URL=https://dtzayqhebbrbvordmabh.supabase.co"
Environment="SUPABASE_BUCKET=spypoint-images"
EnvironmentFile=-/etc/tophand/ranch-camera.env

# Camera Configuration
Environment="DEVICE_NAME=tophand-zero-04"
Environment="IMAGE_WIDTH=2304"
Environment="IMAGE_HEIGHT=1296"
Environment="IMAGE_QUALITY=10"

# Modem Sleep Configuration (saves ~400-500mA)
Environment="MODEM_SLEEP_ENABLED=false"
Environment="MODEM_CONNECTION=hologram"
Environment="MODEM_WAKE_TIME=300"

# Directory Configuration
Environment="IMAGE_DIR=/home/pi/camera/images"
Environment="ARCHIVE_DIR=/home/pi/camera/archive"
Environment="GALLERY_DIR=/home/pi/camera/gallery"

# Restart policy
Restart=on-failure
RestartSec=60

[Install]
WantedBy=multi-user.target
