# ops/systemd/ranchview-dtzay-publisher.timer
[Unit]
Description=Run RanchView dtzay publisher after boot and every 15 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=15min
Persistent=true
RandomizedDelaySec=30s

[Install]
WantedBy=timers.target
